# container/speedtest.crontab
0 * * * *   speedtest --accept-license --accept-gdpr --format json >> /root/speedtest/speedtest.dat
